# control_server_files/roles/webapp/templates/noteapp.service.j2
[Unit]
Description=Note-Taking Web Application
After=network.target
After=httpd.service
Wants=httpd.service

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}/app
Environment=PYTHONPATH={{ app_dir }}/app
Environment=FLASK_APP=app.py
Environment=FLASK_ENV=production
ExecStart=/usr/bin/python3 {{ app_dir }}/app/app.py
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security Settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ app_dir }}
PrivateDevices=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=noteapp

[Install]
WantedBy=multi-user.target 