# control_server_files/templates/simple-note-app.service.j2
[Unit]
Description={{ service_description }}
After=network.target
After=mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ app_dir }}
Environment=PATH={{ virtual_env_path }}/bin
Environment=FLASK_APP=frontend.py
Environment=FLASK_ENV={{ flask_env }}
ExecStart={{ virtual_env_path }}/bin/python {{ app_dir }}/frontend.py
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security Settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ app_base_dir }}
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ service_name }}

[Install]
WantedBy=multi-user.target 